# .config/etc/udev/rules.d/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/neel/.Xauthority", RUN+="/home/neel/.config/monitor-hotplug.sh"
